# systemd/mnt-truenas\x2dlocal-brents\x2ddata.automount
[Unit]
Description=Automount TrueNAS SMB over Local Network

[Automount]
Where=/mnt/truenas-local/brents-data
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target